# crontab
SHELL=/bin/bash
47 16 * * *       /home/paul/dev/ex/git/feature-test/runscript >> /home/paul/dev/ex/git/feature-test/output 2>&1
